# dist/omnomnomnom-web.service
[Unit]
Description=omnomnomnom-web
After=network.target

[Service]
User=omnomnomnom
WorkingDirectory=/home/omnomnomnom/omnomnomnom/api
Environment=FLASK_ENV=production
Environment=APP_SETTINGS='config.production_secret.Config'
ExecStart=/home/omnomnomnom/omnomnomnom/venv/bin/waitress-serve --host 127.0.0.1 --port 5000 --call 'app:create_app'
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
